# supervisor-redash-alert.ini
[program:redash-alert]
command = python3 -m redash_alert.main
process_name = %(program_name)s
directory = /opt/redash-alert
numprocs = 1
startsecs = 5
stopsignal = TERM
stopwaitsecs = 15
autostart = true
autorestart = true
redirect_stderr = true
stopasgroup=true
user = root
group = root
stdout_logfile_maxbytes = 50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups = 10 ; # of stdout logfile backups (default 10)
stdout_logfile = /opt/redash-alert/alert.log
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 10
stderr_logfile = /opt/redash-alert/alert.log
minfds = 65535
minprocs = 32700